# packaging/1985-web.service.in
# 1985-web.service — run the 1985 Amstrad PCW emulator as a sandboxed,
# system-wide systemd service, serving the Web Service to browsers on the
# network (see 1985(1), option --web). --web implies --headless: SDL's
# offscreen video and dummy audio drivers are forced, no window ever
# appears on the host, and the browser is the only interface.
#
# Runs as the dedicated, shell-less 'emulator' system account (shared with
# 1984's equivalent unit, if installed) rather than root or an interactive
# login user, with systemd sandboxing on top of that — see SECURITY below.
#
#   systemctl enable --now 1985-web    # start and enable at boot
#   systemctl stop 1985-web            # stop
#
# The port defaults to web_port in the emulator account's config file
# (/var/lib/emulator/.config/1985/1985.conf). To override it (or add disks,
# pick a model, ...) use a drop-in:
#
#   systemctl edit 1985-web
#   [Service]
#   ExecStart=
#   ExecStart=@BINDIR@/1985 --web=8080 --disk-a=/var/lib/emulator/disks/work.dsk
#
# SECURITY: the Web Service binds 0.0.0.0 with no authentication — anyone
# who can reach the port gets a full, isolated emulator session. Do not
# expose this directly to the internet: put it behind a VPN (Tailscale,
# WireGuard) or a reverse proxy that adds TLS and authentication. Running
# as 'emulator' plus the sandboxing below limits the blast radius if the
# service is ever compromised — it is not a substitute for network-level
# access control.

[Unit]
Description=1985 Amstrad PCW emulator (Web Service, sandboxed)
Documentation=man:1985(1)
After=network.target

[Service]
ExecStart=@BINDIR@/1985 --web
Restart=on-failure
RestartSec=3

User=emulator
Group=emulator
# Creates/owns /var/lib/emulator (also this account's $HOME, per its
# passwd entry) and grants it write access under ProtectSystem=strict.
StateDirectory=emulator

# Offscreen video + dummy audio (forced by --web) need no GPU/audio device
# access; the service needs no elevated privileges and no host filesystem
# access beyond its own state directory.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RemoveIPC=yes
CapabilityBoundingSet=
UMask=0077

[Install]
WantedBy=multi-user.target
