# packaging/1985-web.service.in
# 1985-web.service — run the 1985 Amstrad PCW emulator as a headless
# systemd user service, serving the Web GUI to browsers on the network
# (see 1985(1), option --web). --web implies --headless: SDL's
# offscreen video and dummy audio drivers are forced, no window ever
# appears on the host, and the browser is the only interface.
#
#   systemctl --user enable --now 1985-web    # start and enable
#   systemctl --user stop 1985-web            # stop
#
# The port defaults to web_port in ~/.config/1985/1985.conf (1985).
# To override it (or add disks, pick a model, ...) use a drop-in:
#
#   systemctl --user edit 1985-web
#   [Service]
#   ExecStart=
#   ExecStart=@BINDIR@/1985 --web=8080 --disk-a %h/disks/work.dsk
#
# To start at boot without an active login session:
#   loginctl enable-linger $USER
#
# SECURITY: the Web GUI binds 0.0.0.0 with no authentication — anyone
# on the local network can watch the screen and type into the emulator.

[Unit]
Description=1985 Amstrad PCW emulator (headless Web GUI)
Documentation=man:1985(1)
After=network.target

[Service]
ExecStart=@BINDIR@/1985 --web
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
